# channel_guard.service
[Unit]
Description=LND Channel Guard - Monitors and protects channel outbound liquidity
After=network.target lnd.service

[Service]
User=admin
Group=admin
ExecStart=/usr/local/bin/channel_guard.py CHANGE_ME_CHAN_ID
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"
Environment="HOME=/home/admin"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
